# deploy/systemd/mindmate.service
[Unit]
Description=MindMate Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/.openclaw/workspace-papzin-lead/projects/mindmate
EnvironmentFile=/home/ubuntu/.openclaw/workspace-papzin-lead/.env
EnvironmentFile=/home/ubuntu/.openclaw/workspace-papzin-lead/projects/mindmate/.env
Environment=PYTHONUNBUFFERED=1
Environment=PORT=10000
ExecStart=/home/ubuntu/.openclaw/workspace-papzin-lead/projects/mindmate/venv/bin/python /home/ubuntu/.openclaw/workspace-papzin-lead/projects/mindmate/bot.py
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=default.target
